# custom/supervisord-general.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:unit]
command=/usr/sbin/unitd --no-daemon
autostart=true
autorestart=true
startretries=0
priority=5
stdout_logfile=/var/log/unit/access.log
stderr_logfile=/var/log/unit/error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB