# custom/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:unit]
command=/usr/sbin/unitd --no-daemon
autostart=true
autorestart=true
startretries=0
priority=5

[program:unit-init]
command=/bin/sh -c "until [ -S /run/control.unit.sock ]; do sleep 0.1; done && curl -sf -X PUT --unix-socket /run/control.unit.sock -d '{\"max\":${UNIT_MAX_PROCESSES},\"spare\":${UNIT_SPARE_PROCESSES},\"idle_timeout\":60}' http://localhost/config/applications/default/processes"
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=10

[program:laravel-queue]
command=/usr/bin/php /app/artisan queue:work
autostart=${WITH_QUEUE}
autorestart=true
startretries=0
priority=15

[program:laravel-schedule]
command=/usr/bin/php /app/artisan schedule:work
autostart=${WITH_SCHEDULE}
autorestart=true
startretries=0
priority=20